# wireguard_check_endpoints.timer
[Unit]
Description=Periodically Check Wireguard Endpoint

[Timer]
OnCalendar=OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
